# support/supervisord/config/lb_role_store.ini
[program:lb_role_store]
process_name = %(program_name)s_%(process_num)02d
command = /usr/local/bin/php /Users/xuyakun/data/wwwroot/laravel-9.x-back/artisan queue:work redis --queue=role_store,listeners --daemon --tries=1 --sleep=1
autostart = true
autorestart = true
user = root
numprocs = 1
redirect_stderr = true
stdout_logfile = /usr/local/var/log/supervisor/lb_role_store.log
